# supervisord.conf
[supervisord]
nodaemon=true

[program:sphinx]
command=searchd
autorestart=true